# image/services.d/cron/finish
#!/bin/bash

echo "cron/finish: cron service was stopped"
